# src/NetworkTester.init
#!/bin/sh /etc/rc.common

START=30
SERVICE_DAEMONIZE=1
NAME=NetworkTester

start() {
	service_start /tester/daemon.py
}

stop() {
	service_stop /tester/daemon.py
}

reload() {
	service_reload /tester/daemon.py
}

